# profile-symlinks.list
airootfs/etc/localtime|/usr/share/zoneinfo/UTC
airootfs/etc/resolv.conf|/run/systemd/resolve/stub-resolv.conf
airootfs/etc/systemd/system/cloud-init.target.wants/cloud-config.service|/usr/lib/systemd/system/cloud-config.service
airootfs/etc/systemd/system/cloud-init.target.wants/cloud-final.service|/usr/lib/systemd/system/cloud-final.service
airootfs/etc/systemd/system/cloud-init.target.wants/cloud-init-local.service|/usr/lib/systemd/system/cloud-init-local.service
airootfs/etc/systemd/system/cloud-init.target.wants/cloud-init-main.service|/usr/lib/systemd/system/cloud-init-main.service
airootfs/etc/systemd/system/cloud-init.target.wants/cloud-init-network.service|/usr/lib/systemd/system/cloud-init-network.service
airootfs/etc/systemd/system/dbus-org.freedesktop.ModemManager1.service|/usr/lib/systemd/system/ModemManager.service
airootfs/etc/systemd/system/dbus-org.freedesktop.resolve1.service|/usr/lib/systemd/system/systemd-resolved.service
airootfs/etc/systemd/system/dbus-org.freedesktop.timesync1.service|/usr/lib/systemd/system/systemd-timesyncd.service
airootfs/etc/systemd/system/multi-user.target.wants/ModemManager.service|/usr/lib/systemd/system/ModemManager.service
airootfs/etc/systemd/system/multi-user.target.wants/choose-mirror.service|../choose-mirror.service
airootfs/etc/systemd/system/multi-user.target.wants/hv_fcopy_daemon.service|/usr/lib/systemd/system/hv_fcopy_daemon.service
airootfs/etc/systemd/system/multi-user.target.wants/hv_kvp_daemon.service|/usr/lib/systemd/system/hv_kvp_daemon.service
airootfs/etc/systemd/system/multi-user.target.wants/hv_vss_daemon.service|/usr/lib/systemd/system/hv_vss_daemon.service
airootfs/etc/systemd/system/multi-user.target.wants/livecd-talk.service|/etc/systemd/system/livecd-talk.service
airootfs/etc/systemd/system/multi-user.target.wants/pacman-init.service|../pacman-init.service
airootfs/etc/systemd/system/multi-user.target.wants/systemd-resolved.service|/usr/lib/systemd/system/systemd-resolved.service
airootfs/etc/systemd/system/multi-user.target.wants/vboxservice.service|/usr/lib/systemd/system/vboxservice.service
airootfs/etc/systemd/system/multi-user.target.wants/vmtoolsd.service|/usr/lib/systemd/system/vmtoolsd.service
airootfs/etc/systemd/system/multi-user.target.wants/vmware-vmblock-fuse.service|/usr/lib/systemd/system/vmware-vmblock-fuse.service
airootfs/etc/systemd/system/multi-user.target.wants/NetworkManager.service|/usr/lib/systemd/system/NetworkManager.service
airootfs/etc/systemd/system/network-online.target.wants/NetworkManager-wait-online.service|/usr/lib/systemd/system/NetworkManager-wait-online.service
airootfs/etc/systemd/system/sockets.target.wants/pcscd.socket|/usr/lib/systemd/system/pcscd.socket
airootfs/etc/systemd/system/sound.target.wants/livecd-alsa-unmuter.service|../livecd-alsa-unmuter.service
airootfs/etc/systemd/system/sysinit.target.wants/systemd-time-wait-sync.service|/usr/lib/systemd/system/systemd-time-wait-sync.service
airootfs/etc/systemd/system/sysinit.target.wants/systemd-timesyncd.service|/usr/lib/systemd/system/systemd-timesyncd.service
airootfs/etc/systemd/system/display-manager.service|/usr/lib/systemd/system/sddm.service
airootfs/etc/systemd/system/default.target|/usr/lib/systemd/system/graphical.target
airootfs/etc/systemd/system-generators/systemd-gpt-auto-generator|/dev/null
